# rc.local
#!/bin/sh

# list network interfaces and ip addresses
echo
for IF in $(ls /sys/class/net/ | sed -e '/^lo/d' -e '/^sit/d' | sort); do
	IP=$(ifconfig $IF | grep 'inet addr:' | cut -d: -f2 | awk '{print $1}')
	if [ -n "$IP" ]; then
		echo "$IF: $IP"
	fi
done

exit 0
