# systemd/retroflag.service
[Unit]
Description=retroflag daemon

[Service]
Type=forking
ExecStart=/usr/sbin/retroflag

[Install]
WantedBy=multi-user.target
